# misc/systemd/prosody-c2s.socket
[Install]
WantedBy=sockets.target

[Socket]
ListenStream=0.0.0.0:5222
ListenStream=5222
BindIPv6Only=ipv6-only
Accept=false
Service=prosody.service
